# root/etc/services.d/plex/finish
#!/usr/bin/with-contenv bash
# Copied from the init.d stop method from non-Dockerized Plex.

echo "Stopping Plex Media Server."

# Ask nicely
pids="$(ps -ef | grep 'Plex Media Server' | grep -v grep | awk '{print $2}')"
if [ -n "$pids" ]; then
  kill -15 $pids
fi

# Wait for a graceful exit, but proceed the moment Plex is gone so we keep as
# much of the shutdown window as possible for the backup below.
for _ in $(seq 1 20); do
  pids="$(ps -ef | grep /usr/lib/plexmediaserver | grep -v grep | awk '{print $2}')"
  [ -z "$pids" ] && break
  sleep 0.25
done

# Force-kill anything still alive.
pids="$(ps -ef | grep /usr/lib/plexmediaserver | grep -v grep | awk '{print $2}')"
if [ -n "$pids" ]; then
  kill -9 $pids
  sleep 0.5
fi

# Plex is now fully stopped: back up the databases.
/usr/local/bin/plex-db-backup
